# crontab
* * * * * cd /root && run_crawler.sh >> /var/log/cron.log 2>&1